# src/install/lib/systemd/system/lsc-sync@.timer
[Unit]
Description=Periodically run LSC %I

[Timer]
OnCalendar= *:30:00
