# src/install/lib/systemd/system/lsc-oneshot.timer
[Unit]
Description=Periodically run LSC

[Timer]
OnCalendar= *:30:00
