# lsc-2.1.5/etc/cron.d/lsc.cron
#
# Regular cron jobs for the Ldap Synchronization Connector
#

30 * * * * root [ -x #LSC_BIN# ] && #LSC_BIN# -s all -c all > /dev/null 2>&1
